# build-scope/pkg/general/etc/cron.d/nxs-backup
0 3		* * *	root	/usr/sbin/nxs-backup start all
